# deployment/frontend.service
[Unit]
Description=FlashLab Next.js Frontend
After=network.target

[Service]
User=www-data
WorkingDirectory=/opt/flashlab/frontend
ExecStart=/usr/bin/npm start
Restart=always
Environment=NODE_ENV=production
Environment=NEXT_PUBLIC_API_URL=https://flashlab.pro/api

[Install]
WantedBy=multi-user.target
